# root/etc/init.d/telegrambot
#!/bin/sh /etc/rc.common

START=99
STOP=15

start() {
  /usr/bin/telegrambot.sh &
}

stop() {
  killall telegrambot.sh
}
